# scripts/service/transmission/run
#!/usr/bin/with-contenv sh

/usr/bin/transmission-daemon --foreground --config-dir /etc/transmission-daemon